# crontab
0 */3 * * * /usr/local/bin/rclone sync $HOME/Code storagebox:CodeBackup

